# cpr_telegram.service
[Unit]
Description=CPR Telegram Bot Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/cpr_bot
ExecStart=/home/ubuntu/cpr_bot/venv/bin/python3 /home/ubuntu/cpr_bot/telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/cpr_telegram.log
StandardError=append:/var/log/cpr_telegram.err
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
